# Build/LinUtil/files/prl-disp.service
[Unit]
Description=Dispatcher service.
After=syslog.target network.target
After=libvirtd.service systemd-modules-load.service vzevent.service vcmmd.service vz.service
Wants=libvirtd.service systemd-modules-load.service vzevent.service vcmmd.service vz.service

[Service]
Type=forking
PIDFile=/var/run/prl_disp_service.pid
ExecStart=/usr/sbin/prl_disp_service --pidfile /var/run/prl_disp_service.pid
KillMode=process
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
