# Build/LinUtil/files/virtuozzo.target
[Unit]
Description=Virtuozzo virtual machines and container target
PartOf=prl-disp.service
