# Build/LinUtil/files/virtuozzo.target
[Unit]
Description=Virtuozzo virtual machines and container target
BindsTo=prl-disp.service
ConditionPathExists=/run/virtuozzo.target.start
